# deploy/hardening/cloudflare-real-ip-update.timer
[Unit]
Description=Refresh Cloudflare proxy ranges weekly

[Timer]
OnCalendar=weekly
RandomizedDelaySec=6h
Persistent=true

[Install]
WantedBy=timers.target
